# demo/RPi_Air/rpi_service_air.service
[Unit]
Description=rpi_service

[Service]
ExecStart=/root/RPi_Air/rpi_sh.sh

[Install]
WantedBy=multi-user.target
