# demo/rpi/rpi_service.service
[Unit]
Description=rpi_service

[Service]
ExecStart=/root/RPi/rpi_sh.sh

[Install]
WantedBy=multi-user.target
